# system/auris-record.service
[Unit]
Description=Auris audio recorder (ALSA to file)

[Service]
Type=simple
ExecStart=/opt/auris/record.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
